# relay/crontab
* * * * * /vlab/checkboards.py -s -k -v >> /vlab/log/relay.log 2>&1

