# relay/crontab
* * * * * /vlab/checkboards.py -s -k -v >> /var/log/vlab.log 2>&1

