# relay/crontab
* * * * * /vlab/checkboards.py -s -k >> /vlab/log/relay.log 2>&1

