# host/etc/cron.d/vlab-boardscan
* * * * * root /opt/VLAB/boardscan.sh >> /opt/VLAB/log/boardscan.log 2>&1
