# boardserver/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:hwserver]
command=/opt/Xilinx/HWSRVR/current/bin/hw_server
autorestart=true

[program:cron]
command=cron -f
