# boardserver/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:hwserver]
command=/opt/Xilinx/HWSRVR/2016.2/bin/hw_server

[program:cron]
command=cron -f
